# extras/xenial64.extras-tty.list
etc/
etc/rc.local
etc/systemd/
etc/systemd/system/
etc/systemd/system/rclocal-shutdown.service
etc/rc.local.shutdown
etc/console-setup/
etc/console-setup/mmgen-24.psfu
home/
home/mmgen/
home/mmgen/.bash_profile
home/mmgen/README.fullnode
home/mmgen/scripts/
home/mmgen/scripts/term/
home/mmgen/scripts/term/vars.sh
home/mmgen/scripts/term/terminal.sh
home/mmgen/scripts/term/terminal-small.sh
home/mmgen/scripts/term/terminal-large.sh
home/mmgen/scripts/term/tmux.sh
home/mmgen/scripts/term/terminal-medium.sh
home/mmgen/scripts/include/
home/mmgen/scripts/include/functions.sh
home/mmgen/doc/
home/mmgen/doc/Getting-Started-with-MMGen
home/mmgen/doc/README.MMGenLive
home/mmgen/doc/README.Unix
home/mmgen/bin/
home/mmgen/bin/mmlive-node-setup
home/mmgen/bin/mmlive-wifi-disable
home/mmgen/bin/mmlive-disk-passwd-add
home/mmgen/bin/mmlive-disk-passwd-remove
home/mmgen/bin/mmlive-networking-disable
home/mmgen/bin/mmlive-node-start
home/mmgen/bin/mmlive-bluetooth-disable
home/mmgen/bin/mmlive-bluetooth-enable
home/mmgen/bin/mmgen-test
home/mmgen/bin/mmlive-networking-enable
home/mmgen/bin/mmlive-onion-nodes-add
home/mmgen/bin/mmlive-autologin-disable
home/mmgen/bin/mmlive-autologin-enable
home/mmgen/bin/mmlive-wifi-enable
home/mmgen/bin/mmlive-user-passwd-chg
home/mmgen/README
home/mmgen/.tmux.conf
home/mmgen/.bashrc
home/mmgen/.vimrc
root/
root/.bash_login
usr/
usr/local/
usr/local/sbin/
usr/local/sbin/update-grub-cfg-kver.sh
