# extras/xenial64.extras-tty.list
etc/
etc/rc.local
etc/systemd/
etc/systemd/system/
etc/systemd/system/rclocal-shutdown.service
etc/rc.local.shutdown
etc/console-setup/
etc/console-setup/mmgen-24.psfu
home/
home/mmgen/
home/mmgen/.bash_profile
home/mmgen/.tmux.conf
home/mmgen/.bashrc
home/mmgen/.vimrc
home/mmgen/.gvimrc
root/
root/.bash_login
root/.vimrc
root/.vim
usr/
usr/local/
usr/local/sbin/
usr/local/sbin/update-grub-cfg-kver.sh
